# tor/service/privoxy/run
#!/bin/sh
privoxy --no-daemon /root/configs/privoxy/config
