# tor/service/tor/run
#!/bin/sh
tor -f /root/configs/tor/torrc
